# plasma-kwin_wayland.service.in
[Unit]
Description=KDE Window Manager
PartOf=graphical-session.target

[Service]
Type=notify-reload
Restart=on-failure
FileDescriptorStoreMax=200
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_wayland --xwayland
Slice=session.slice
WatchdogSec=15s
NotifyAccess=all
WatchdogSignal=SIGHUP
